# contrib/solstone-tmux.service
[Unit]
Description=Solstone Tmux Terminal Observer
After=basic.target

[Service]
Type=simple
# Environment=PATH=... (injected at install time by solstone-tmux install-service)
ExecStart=%h/.local/bin/solstone-tmux run
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target
