# src/solstone_tmux/contrib/solstone-tmux.service.in
[Unit]
Description=Solstone Tmux Terminal Observer
After=basic.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
Environment=PATH={PATH}
ExecStart={BINARY} run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
